# packaging/com.yourorg.codexpc.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.yourorg.codexpc</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/codexpc/bin/codexpcd</string>
  </array>
  <key>MachServices</key>
  <dict>
    <key>com.yourorg.codexpc</key>
    <true/>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>$HOME/Library/Logs/com.yourorg.codexpc.out.log</string>
  <key>StandardErrorPath</key>
  <string>$HOME/Library/Logs/com.yourorg.codexpc.err.log</string>
</dict>
</plist>
